# com.claude-copy.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.claude-copy</string>

  <key>ProgramArguments</key>
  <array>
    <string>__PYTHON_PATH__</string>
    <string>-u</string>
    <string>__DAEMON_PATH__</string>
  </array>

  <key>SessionCreate</key>
  <true/>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>30</integer>

  <key>StandardOutPath</key>
  <string>/tmp/claude-copy.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/claude-copy.log</string>
</dict>
</plist>
